# setups/docker/docker-system-prune.timer
[Unit]
Description=Run docker system prune weekly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target